# dnsmasq-adblock/etc/cron.weekly/dnsmasq-update-adblock-hosts
#!/bin/sh

cd /usr/local/etc/dnsmasq/adblock
./generate || exit
#systemctl reload dnsmasq

sudo -u shared-config ~shared-config/bin/push-dnsmasq "dnsmasq-adblock-update-hosts $(date -Iseconds)"
