# dnsmasq-adblock/etc/cron.weekly/dnsmasq-update-adblock-hosts
#!/bin/bash

cd /var/lib/dnsmasq/adblock
./generate
systemctl reload dnsmasq
